# assets/supervisord.conf
[supervisord]
nodaemon=true

[program:shadowsocks]
command=/usr/bin/ss-local -c /etc/shadowsocks-libev/config.json -v
