# assets/supervisord.conf
[supervisord]
nodaemon=true

[program:shadowsocks]
command=/usr/sbin/ss-server -c /etc/shadowsocks-libev/config.json -v
